# root/etc/s6-overlay/s6-rc.d/init-mod-universal-wireguardvpn-add-package/run
#!/usr/bin/with-contenv bash

## Ubuntu
if [ -f /usr/bin/apt ]; then
    echo "**** Adding universal-wireguardvpn packages to install lists ****"
    echo "\
        systemd-resolved \
        wireguard-tools" >> /mod-repo-packages-to-install.list

fi

# Alpine
if [ -f /sbin/apk ]; then
    if ! apk info --installed wireguard-tools >/dev/null; then
        echo "**** Adding universal-wireguardvpn packages to install lists ****"
        echo "\
            grep \
            iproute2 \
            iptables \
            iptables-legacy \
            ip6tables \
            iputils \
            libcap-utils \
            net-tools \
            wireguard-tools" >> /mod-repo-packages-to-install.list
    else
        echo "**** Universal-wireguardvpn packages already installed, skipping ****"
    fi
fi
