# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-wireguardvpn/run
#!/usr/bin/with-contenv bash

echo "************ Wireguard VPN **************"
wg-quick
sleep infinity
